# 10-Huawei.rules
SUBSYSTEM=="tty", ATTRS{idProduct}=="1506", ATTRS{idVendor}=="12d1", RUN+="/home/pi/rtsp/misc/dongle/connect.sh"
